# deploy/wendio.service
[Unit]
Description=WENDIO — Bingo Wendat Socket.io
Documentation=https://github.com/jeffze/wendio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=darkvador
Group=darkvador
WorkingDirectory=/home/darkvador/jeux/wendio
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3

# Bind loopback : Caddy gère l'exposition publique
Environment=PORT=5000
Environment=HOST=127.0.0.1
Environment=NODE_ENV=production

# Logs vers journald (consulter avec : sudo journalctl -u wendio -f)
StandardOutput=journal
StandardError=journal

# Durcissement basique
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/darkvador/jeux/wendio
PrivateTmp=true
PrivateDevices=true

# Limites raisonnables
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
